# app/frontend/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/ft_transcendence/public/logs/supervisord.log
pidfile=/home/ft_transcendence/public/logs/supervisord.pid
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
